# root/etc/crontabs/root
*/5 * * * * sleep $((0 + $RANDOM % 120)); /app/duck.sh