# root/etc/crontabs/root
*/5 * * * * sleep $(0 + $RANDOM % 240); /app/duck.sh